# packages/llmd-engine-ollama/lib/systemd/system/ollama.service.d/hardening.conf
# SibillaOS: hardening drop-in for the Ollama service
[Service]
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/llmd/models /usr/share/ollama
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
# Models live in the shared data partition
Environment=OLLAMA_MODELS=/var/lib/llmd/models/ollama
# Loopback only: external exposure goes through the gateway
Environment=OLLAMA_HOST=127.0.0.1:11434
