# packages/llmd-mcp/lib/systemd/system/llmd-mcp.service
[Unit]
Description=SibillaOS Model Context Protocol server
Documentation=https://github.com/engineering87/sibillaos
After=network.target

[Service]
# notify: startup completes when the socket is bound, so
# `systemctl start` (and sibilla mcp enable with it) returns only
# once the server actually answers - no race for the first request
Type=notify
ExecStart=/usr/bin/python3 /usr/lib/llmd/llmd-mcp-server
Restart=on-failure
RestartSec=2
# The server reads only /etc/llmd/engine and /etc/llmd/model, never
# the API keys: authentication lives at the gateway. A dynamic user
# with the standard SibillaOS sandbox is enough.
DynamicUser=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
