# packages/llmd-gateway/lib/systemd/system/caddy.service.d/apikey.conf
# SibillaOS: pass the generated API key to Caddy
# /etc/llmd/gateway.env is written by llmd-firstboot (mode 0600)
[Service]
EnvironmentFile=-/etc/llmd/gateway.env
