# packages/llmd-gateway/lib/systemd/system/caddy.service.d/hardening.conf
# SibillaOS: hardening drop-in for the Caddy gateway.
# CAP_NET_BIND_SERVICE stays available through the ambient capabilities
# of the stock unit: ambient caps survive NoNewPrivileges by design.
# No ProtectSystem=strict here: the stock unit already sets full, and
# Caddy needs its state directory (/var/lib/caddy) for the local CA.
[Service]
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
