# packages/llmd-engine-vllm/etc/containers/systemd/vllm.container
# SibillaOS Quadlet: vLLM in an OCI container via podman
# Generates vllm.service automatically (systemd >= 253, podman >= 4.4)
[Unit]
Description=vLLM OpenAI-compatible server (container)
After=network-online.target
Wants=network-online.target
# Only activate on machines where the installer selected vLLM:
# /etc/llmd/vllm.env is written by llmd-firstboot in that case only.
# Without this gate the unit would start (and pull a multi-GB image)
# on every machine, Ollama ones included.
ConditionPathExists=/etc/llmd/vllm.env

[Container]
# PINNED, like every engine in this project: ":latest" survived here
# from the PoC and was found during the v0.7 supply-freshness pass.
# Bump deliberately; the quadlet dry-run in lint and the vllm gating
# checks in the install test cover the unit, the model serving path
# stays validated on real GPUs only (see the roadmap).
Image=docker.io/vllm/vllm-openai:v0.25.1
# NVIDIA GPU via CDI (nvidia-container-toolkit)
AddDevice=nvidia.com/gpu=all
PublishPort=127.0.0.1:8000:8000
Volume=/var/lib/llmd/models/hf:/root/.cache/huggingface:Z
EnvironmentFile=-/etc/llmd/vllm.env
# the workload never needs to gain privileges at runtime
PodmanArgs=--security-opt=no-new-privileges
# The model selected at install time is passed as an argument
Exec=--model ${LLMD_MODEL} --host 0.0.0.0 --port 8000

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
