# packages/llmd-firstboot/lib/systemd/system/llmd-firstboot.service
[Unit]
Description=SibillaOS first boot: model download and service setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/llmd/firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/lib/llmd/llmd-firstboot
RemainAfterExit=true
TimeoutStartSec=infinity
StandardOutput=journal+console
StandardError=journal+console
# Sandbox, as far as a root setup task allows. No ProtectHome: the
# ollama CLI writes its key under /root/.ollama during the pull.
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
